# systemd-units/update-container-sites.service
[Unit]
Description=Обновление реестра контейнерных площадок

[Service]
Type=oneshot
ExecStart=%h/trash/tools/update-container-sites.sh
